# files/pgpool-install.service
[Unit]
Description=Install pgpool

[Service]
Type=simple
User=root

ExecStart=/usr/bin/pgpool-install

[Install]
WantedBy=multi-user.target
