# files/pgpool-install.service
[Unit]
Description=Install pgpool

[Service]
Type=simple
User=root

ExecStart=/usr/bin/postgresql_create_root

[Install]
WantedBy=multi-user.target
