# files/pgpool_node_create.service
[Unit]
Description=Pgpool node

[Service]
Type=simple
User=root

ExecStart=/usr/bin/pgpool_node_create

[Install]
WantedBy=multi-user.target

